# systemd/system/multi-user.target.wants/weather-rawlog.service
[Unit]
Description=Weatherstation: rawlog
After=remote-fs.target

[Service]
ExecStart=/usr/local/lib/weather/rawlog.py
Restart=always
