# systemd/system/multi-user.target.wants/weather.service
[Unit]
Description=Weatherstation: relog
After=remote-fs.target weather-rawlog.service

[Service]
ExecStart=/usr/local/lib/weather/relog.py
Restart=always
